# bizmark-api.service
[Unit]
Description=Bizmark.id SaaS Platform Perizinan UMKM
Documentation=https://github.com/tanerizawa/bizmark.id
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/var/www/bizmark.id/backend
ExecStart=/usr/bin/npm run start:prod
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=bizmark-api
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
